# rootfs/etc/services.d/php-fpm8.0/run
#!/usr/bin/with-contenv bash
echo "Starting PHP-FPM";
exec /usr/sbin/php-fpm8.0 -F

